# pkg/xterm/etc/services.d/xterm/run
#!/usr/bin/with-contenv bash

sudo -Hu user xterm