# pkg/websocketify/etc/services.d/websocketify/run
#!/usr/bin/with-contenv bash
websockify --web=/usr/share/novnc/ --cert=/etc/ssl/novnc.pem 80 127.0.0.1:5900