# pkg-wine/wine/etc/services.d/wine/run
#!/usr/bin/with-contenv bash

set -eo pipefail
[[ $DEBUG == true ]] && set -x

tail /dev/null -F

sleep 1

wine wineboot --kill
wineserver -k9