# pkg-wine/wine/etc/services.d/wine/run
#!/usr/bin/with-contenv bash

set -eo pipefail
[[ $DEBUG == true ]] && set -x

wine explorer /desktop=windows,${VNC_GEOMETRY} notepad

sleep 1

wine wineboot --kill
wineserver -k9